# scripts/cdrommon@.service
[Unit]
Description=Monitor CDROM of %I

[Service]
ExecStart=/opt/xensource/libexec/cdrommon /dev/xapi/cd/%I
