# scripts/mpathalert-daemon.service
[Unit]
Description=MPathAlert Daemon Service

[Service]
ExecStart=@BINDIR@/mpathalert -delay 120
StandardError=null
