# scripts/wsproxy.socket
[Socket]
ListenStream=/var/lib/xcp/websockproxy

[Install]
WantedBy=sockets.target
