# scripts/wsproxy.service
[Unit]
Description=Websocket proxy

[Service]
ExecStart=@LIBEXECDIR@/wsproxy
StandardInput=socket

[Install]
WantedBy=multi-user.target
