# scripts/wsproxy.socket
[Socket]
ListenStream=/var/xapi/wsproxy

[Install]
WantedBy=sockets.target
