# scripts/xapi-init-complete.target
[Unit]
Description=Xapi is available
After=syslog.target
After=time-sync.target
