# scripts/xen-frontend.rules
SUBSYSTEM=="block", KERNEL=="xvd*", RUN+="@SCRIPTSDIR@/block-frontend $env{ACTION} $number"
