# scripts/xen-frontend.rules
SUBSYSTEM=="block", KERNEL=="xvd*", RUN+="/etc/xensource/scripts/block-frontend $env{ACTION} $number"
