# vendor/xapi-xenopsd/scripts/xen-backend.rules.in
SUBSYSTEM=="xen-backend", KERNEL=="tap*", RUN+="@LIBEXEC@/tap $env{ACTION}"
SUBSYSTEM=="xen-backend", KERNEL=="vbd*|qdisk", RUN+="@LIBEXEC@/block $env{ACTION}"

SUBSYSTEM=="xen-backend", KERNEL=="vif*", RUN+="@LIBEXEC@/vif $env{ACTION} type_if=vif"
SUBSYSTEM=="net",         KERNEL=="tap*", RUN+="@LIBEXEC@/vif $env{ACTION} type_if=tap"

KERNEL=="evtchn", NAME="xen/%k"
